# myWork/dependencies/human.service
[Unit]
Description=This is human service

[Service]
ExecStart=echo unga bunga
RemainAfterExit=yes

