# myWork/dependencies/water.service
[Unit]
description=This is water service

[Service]
ExecStart=/bin/false
